# rclone/services.d/rclone/run
#!/usr/bin/with-contenv bash

# UMASK_SET=${UMASK_SET:-022}
# umask "$UMASK_SET"

echo "Starting Rclone."
exec \
        /bin/bash -c \
        '/usr/bin/mount-rclone'
